# overlay/generic/etc/cron.d/crontabs/root
## Rotate and trim the audit logs nightly.
0 0 * * * /smartdc/bin/sdc-lastcomm -R 30
## Run logadm hourly so that the minimum rotate interval is 1h
0 * * * * /usr/sbin/logadm
## Delete saved core dumps over 7 days old
15 0 * * * find /zones/*/cores -type f -mtime +7 -exec rm -f "{}" \;

# NOTE: all entries above are in root's system-defined crontab; see the
# crontab(1) man page.

