# overlay/var/spool/cron/crontabs/root
## Rotate and trim the audit logs nightly at midnight.
0 0 * * * /usr/sbin/joycomm -R 30
