# src/vm/smf/system-vmadmd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="vmadmd">
  <service name="system/smartdc/vmadmd" type="service" version="0.0.1">

    <create_default_instance enabled="true"/>
    <single_instance/>

    <!-- smartdc-init also depends on zones, so we can be sure zones is ready too -->
    <dependency name="smartdc-init" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/smartdc/init"/>
    </dependency>

    <dependency name="vminfod" grouping="require_all" restart_on="none" type="service" >
      <service_fmri value="svc:/system/smartdc/vminfod:default" />
    </dependency>

    <exec_method
      type="method"
      name="start"
      exec="/usr/vm/smf/system-vmadmd %m"
      timeout_seconds="60">
      <method_context>
        <method_credential user="root" group="staff"/>
      </method_context>
    </exec_method>

    <exec_method type="method" name="restart" exec=":kill" timeout_seconds="60">
      <method_context>
        <method_credential user="root" group="staff"/>
      </method_context>
    </exec_method>

    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60">
      <method_context>
        <method_credential user="root" group="staff"/>
      </method_context>
    </exec_method>

    <property_group name="startd" type="framework">
      <propval name="ignore_error" type="astring" value="core,signal"/>
    </property_group>

    <property_group name="application" type="application">
    </property_group>

    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">VM Management Daemon (node)</loctext>
      </common_name>
    </template>

  </service>

</service_bundle>
